# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch -Des.config=/etc/elasticsearch/elasticsearch.yml
autorestart=true
